# systemd/mmp-go-reload@.timer
[Timer]
OnActiveSec=1d
OnUnitActiveSec=1d
RandomizedDelaySec=30min

[Install]
WantedBy=timer.target
